# scripts/com.elisa.sync-jira-time.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.elisa.sync-jira-time</string>

    <!-- Paths are injected by install-jira-time-cron.sh at symlink time. -->
    <key>ProgramArguments</key>
    <array>
        <string>__NODE_PATH__</string>
        <string>__SCRIPT_PATH__</string>
    </array>

    <!-- Run at 23:00 local time every day. If the laptop is asleep at 23:00,
         launchd fires the job as soon as it next wakes. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>23</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>__LOG_PATH__</string>
    <key>StandardErrorPath</key>
    <string>__LOG_PATH__</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
