# scripts/com.elisa.sync-airtable-courses.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.elisa.sync-airtable-courses</string>

    <!-- Paths are injected by install-airtable-cron.sh at symlink time. -->
    <key>ProgramArguments</key>
    <array>
        <string>__NODE_PATH__</string>
        <string>__SCRIPT_PATH__</string>
    </array>

    <!-- Run at 22:30 local time every day (30 min before the Jira time sync at
         23:00). If the laptop is asleep, launchd fires the job on next wake. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>22</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>__LOG_PATH__</string>
    <key>StandardErrorPath</key>
    <string>__LOG_PATH__</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
